# nas_root/etc/systemd/system/cloudyhome-nas-apply.service
[Unit]
Description=Apply NAS configuration (datasets, users, services)
After=cloudyhome-nas-render.service cloudyhome-nas-firewall.service cloudyhome-zfs-import.service
Requires=cloudyhome-nas-render.service cloudyhome-nas-firewall.service cloudyhome-zfs-import.service
Wants=cloudyhome-nas-ready.target

[Service]
Type=oneshot
RemainAfterExit=yes
RuntimeDirectory=nas
Environment=SOPS_AGE_KEY_FILE=/etc/sops/age/keys.txt
ExecStart=/usr/local/sbin/nas-apply-config

[Install]
WantedBy=multi-user.target
