# nas_root/etc/systemd/system/cloudyhome-nas-ready.target
[Unit]
Description=NAS provisioning complete
After=cloudyhome-nas-apply.service nfs-server.service smbd.service target.service cloudyhome-garage.service cloudyhome-garage-bootstrap.service cloudyhome-ftp.service
