# nas_root/etc/cloudyhome/templates/cloudyhome-zfs-scrub.service.j2
[Unit]
Description=ZFS pool scrub
After=cloudyhome-zfs-import.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/zpool scrub {{ config.storage.pool }}
